# service/boneserver.service
[Unit]
Description=bonescript web API
 
[Service]
ExecStart=/bin/sh /opt/boneserver/service/boneserver.sh
Restart=on-failure
